# ringboard.slice
[Unit]
Description=Slice for all Ringboard services.
